# packages/ginkos-defaults/etc/update-motd.d/00-ginkos
#!/bin/bash
echo ""
echo "  ░██████╗░██╗███╗░░██╗██╗░░██╗░█████╗░░██████╗"
echo "  ██╔════╝░██║████╗░██║██║░██╔╝██╔══██╗██╔════╝"
echo "  ██║░░██╗░██║██╔██╗██║█████═╝░██║░░██║╚█████╗░"
echo "  ██║░░╚██╗██║██║╚████║██╔═██╗░██║░░██║░╚═══██╗"
echo "  ╚██████╔╝██║██║░╚███║██║░╚██╗╚█████╔╝██████╔╝"
echo "  ░╚═════╝░╚═╝╚═╝░░╚══╝╚═╝░░╚═╝░╚════╝░╚═════╝░"
echo ""
echo "  GinkOS 1.0 (maidenhair) — rooted in the kernel, branching outward"
echo ""
